# roles/webdb/templates/takahe.service.j2
[Unit]
Description=Takahe
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory={{ takahe_checkout_path }}
ExecStart=docker-compose --env-file {{ takahe_env_file }} up
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
